# adsb-multi-portal-mqtt-sensors-fr24/rootfs/etc/s6-overlay/s6-rc.d/adsb-mqtt/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start ha adsb mqtt service
# ==============================================================================

set -eo pipefail

/usr/local/share/adsb-mqtt-sensors/fr24 2>&1
